# etc/service/magpay-api/run
#!/bin/sh
cd /usr/app/api
exec yarn start