# etc/service/magpay-client/run
#!/bin/sh
export PORT=3012
cd /usr/app/client
exec yarn start:live