# etc/service/magpay-client/run
#!/bin/sh
export PORT=3012
exec /usr/app/client yarn build